# pkg/deploy/embedded/clonr-verify-boot.service
[Unit]
Description=clonr post-reboot verification phone-home
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/clonr/node-token
ConditionPathExists=/etc/clonr/verify-boot-url

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clonr-verify-boot
RemainAfterExit=yes
TimeoutStartSec=60
# Retry on network failures
Restart=on-failure
RestartSec=30
StartLimitBurst=10
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
